# sberbot.timer
# должен лежать в /etc/systemd/system/sberbot.timer
# что это и как работает: about_timer.md

[Unit]
Description=Monthly reminder to pick SberSpasibo cashback categories
Requires=sberbot.service

[Timer]
OnCalendar=*-*-01 09:00:00
TimeZone=Europe/Moscow
Persistent=true

[Install]
WantedBy=timers.target
